# deploy/systemd/psfn-face-tracker.service
[Unit]
Description=PSFN Native Face Tracker
After=network-online.target psfn-amica.service
Wants=network-online.target psfn-amica.service

[Service]
Type=simple
User=psfn
WorkingDirectory=/home/psfn/psfn-amica
EnvironmentFile=/home/psfn/psfn-amica/.env.local
Environment=AMICA_FACE_TRACKING_URL=http://127.0.0.1:3000/api/satelliteBridge/
Environment=AMICA_FACE_TRACKING_CAMERA_INDEX=0
ExecStart=/usr/bin/python3 /home/psfn/psfn-amica/scripts/pi_face_tracker.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
